# arp_logger.service
[Unit]
Description=ARP Logger
After=network.target

[Service]
Type=simple
ExecStart=/path/to/your/script.sh
Restart=on-failure
User=root
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
